# hosts/inky/files/etc/systemd/system/inkypi.service.d/10-env.conf
# Load Immich API key and other plugin secrets
[Service]
EnvironmentFile=/etc/inkypi.env
