# hosts/inky/files/etc/systemd/system/mpd.service.d/10-smb.conf
# MPD depends on SMB media mount being available
[Unit]
RequiresMountsFor=/mnt/homelab-media

[Service]
Restart=on-failure
RestartSec=10s
